# infra/stream-host/systemd/blob-chromium.service
[Unit]
Description=Blob stream — Chromium rendering the Stream page
After=blob-xvfb.service blob-pulse.service
Requires=blob-xvfb.service
# Wants, not Requires: if PulseAudio is down the render must still run — it just
# loses sound, which ffmpeg degrades to music-only. A missing sink is not a
# reason to take the picture off the air.
Wants=blob-pulse.service

# A crash loop here must not become a permanent stop: this is the unit the
# watchdog restarts, and if systemd's default start limit (5 starts / 10s) trips
# during a rough patch, the watchdog's reloads silently stop working and the
# stream freezes for good. See blob-ffmpeg.service for why this must be in
# [Unit] and not [Service].
StartLimitIntervalSec=0

[Service]
Type=simple

# NOT root. Chromium hard-refuses to start as root ("Running as root without
# --no-sandbox is not supported"), and on Ubuntu — where `chromium` is really
# the snap — snap confinement additionally restricts the browser to this user's
# $HOME, which is why the profile lives there and not under /run. Running the
# browser unprivileged is the correct posture regardless.
User=blob
Group=blob
Environment=HOME=/home/blob

EnvironmentFile=/opt/blob-stream/.env
ExecStartPre=/bin/sleep 3
ExecStart=/opt/blob-stream/chromium.sh
Restart=always
RestartSec=5

# FENCE THE RENDER OFF THE ENCODER'S CORES.
# blob-ffmpeg is pinned to 0-1; without a complementary pin here Chromium roamed
# all four and fought it. Measured on a live broadcast: the `chrome
# --type=gpu-process` alone holds ~97% of a core doing SOFTWARE compositing (it
# spawns and composites even with --disable-gpu — the display compositor lives
# there regardless), so unfenced it lands on top of the encoder and the picture
# lags while `speed` still reads 1.0x.
#
# 2-3 is the complement of the encoder's 0-1: two cores each, no overlap.
CPUAffinity=2-3

# This is the unit the watchdog restarts. Reloading the BROWSER is the whole
# point: when Streamlit drops its websocket the page freezes while ffmpeg
# happily encodes the dead frame, and restarting the encoder against a frozen
# Chromium just re-encodes the same still image.
[Install]
WantedBy=multi-user.target
